# capistrano/templates/container.service.erb
[Unit]
Description=<%= description %>
After=docker.service
Requires=docker.service

[Service]
ExecStartPre=-/usr/bin/docker kill <%= container %>
ExecStartPre=-/usr/bin/docker rm <%= container %>
ExecStart=/usr/bin/docker run <%= options %> <%= image %> <%= command %>

[Install]
WantedBy=multi-user.target
